# mcp-server.service
[Unit]
Description=MCP Server Control - AI-powered server management
Documentation=https://modelcontextprotocol.io
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/mcp-server
EnvironmentFile=/root/mcp-server/.env
ExecStart=/usr/bin/node /root/mcp-server/server.js
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-server

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
