# deploy/mcp-sentinel-audit-verify.timer
[Unit]
Description=Daily MCP Sentinel audit-chain verification

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
