# deploy/mcp-sentinel-audit-verify.service
[Unit]
Description=Verify MCP Sentinel HMAC audit chain

[Service]
Type=oneshot
User=mcp-sentinel
Group=mcp-sentinel
WorkingDirectory=/opt/mcp-sentinel/current
ExecStart=/usr/bin/node /opt/mcp-sentinel/current/scripts/verify-audit.js
Environment=AUDIT_VERIFICATION_STATUS_FILE=/var/lib/mcp-sentinel/audit-verification.json
EnvironmentFile=/etc/mcp-sentinel/environment
LoadCredential=audit-key:/etc/mcp-sentinel/credentials/audit-key
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/log/mcp-sentinel
ReadWritePaths=/var/lib/mcp-sentinel
