# deploy/mcp-sentinel-broker.service
[Unit]
Description=MCP Sentinel typed privilege broker
Before=mcp-sentinel.service

[Service]
Type=simple
User=root
Group=mcp-sentinel
WorkingDirectory=/opt/mcp-sentinel/current
ExecStart=/usr/bin/node /opt/mcp-sentinel/current/broker.js
EnvironmentFile=/etc/mcp-sentinel/broker-environment
UMask=0077
RuntimeDirectory=mcp-sentinel
RuntimeDirectoryMode=0750
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/run/mcp-sentinel /var/lib/mcp-sentinel /etc/ufw /etc/mcp-sentinel /etc/nginx /etc/passwd /etc/shadow /etc/group /etc/gshadow
# Exact registered project write paths are installed as a root-owned generated
# drop-in by the signed release activator.
Restart=on-failure

[Install]
WantedBy=multi-user.target
